# deb/template.postinst
#!/bin/sh
set -e

if command -v udevadm >/dev/null 2>&1; then
    udevadm control --reload-rules || true
    udevadm trigger || true
fi

echo "Installed @NAME@."
echo "Enable/start per-user service with: systemctl --user enable --now @NAME@.service"
